# scripts/deployment/intent-test-framework.service
[Unit]
Description=AI4SE工具集 - Intent Test Framework
After=network.target postgresql.service

[Service]
Type=simple
User=appuser
Group=appuser
WorkingDirectory=/opt/intent-test-framework
Environment="PATH=/opt/intent-test-framework/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/opt/intent-test-framework/.env
ExecStart=/opt/intent-test-framework/venv/bin/python /opt/intent-test-framework/scripts/deployment/start-production.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=intent-test-framework

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/intent-test-framework

[Install]
WantedBy=multi-user.target


